# .crontab
*/10 * * * * /Users/Hans/icecube/stuff/2014/polecam/script.sh
